# rootfs/etc/services.d/fcgiwrap/run
#!/usr/bin/execlineb -P
su www-data --shell /bin/bash --command "fcgiwrap -s unix:/run/php/fcgiwrap.sock -c 10"
